# templates/vector.service.j2
[Unit]
Description=Vector
Requires=network-online.target
After=network-online.target

[Service]
User=vector
Group=vector
ExecStart=/usr/bin/vector
Restart=always

[Install]
WantedBy=multi-user.target